# scripts/monitoring/saraiva-health-check.service
[Unit]
Description=Saraiva Vision Daily Health Check
After=network.target nginx.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/saraiva-vision-site
ExecStart=/bin/bash /home/saraiva-vision-site/scripts/monitoring/daily-health-check.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=saraiva-health-check

# Timeout de 5 minutos (300 segundos)
TimeoutStartSec=300

# Restart on failure
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
